# rradio.service
[Unit]
Description="RRadio service"

[Service]
ExecStart=/home/pi/rradio/rradio.sh

[Install]
WantedBy=multi-user.target
